# archive/kopia/crontab
0 0 * * * /kopia-scripts/kopia-sync-cron
